# deb/openmediavault/etc/systemd/system/collectd.service.d/openmediavault.conf
[Unit]
# The rrdcached plugin needs a running rrdcached service.
After=rrdcached.service
Wants=rrdcached.service
