# deb/openmediavault/etc/systemd/system/sec.service.d/openmediavault.conf
[Service]
Type=simple
EnvironmentFile=
ExecStart=
ExecStart=/usr/bin/bash -c "/usr/bin/journalctl -f | /usr/bin/sec --conf=/etc/sec.d/*.sec --input=- --notail"
