# deb/openmediavault/etc/cron.weekly/openmediavault-update-smart-drivedb
#!/bin/sh
#
# This file is part of OpenMediaVault.
#
# @license   http://www.gnu.org/licenses/gpl.html GPL Version 3
# @author    Volker Theile <volker.theile@openmediavault.org>
# @copyright Copyright (c) 2009-2015 Volker Theile
#
# OpenMediaVault is free software: you can redistribute it and/or modify
# it under the terms of the GNU General Public License as published by
# the Free Software Foundation, either version 3 of the License, or
# any later version.
#
# OpenMediaVault is distributed in the hope that it will be useful,
# but WITHOUT ANY WARRANTY; without even the implied warranty of
# MERCHANTABILITY or FITNESS FOR A PARTICULAR PURPOSE.  See the
# GNU General Public License for more details.
#
# You should have received a copy of the GNU General Public License
# along with OpenMediaVault. If not, see <http://www.gnu.org/licenses/>.

# Documentation/Howto:
# http://www.thomas-krenn.com/de/wiki/Smartmontools_Drive_Database_aktualisieren

set -e

tmpscript=$(mktemp)

cleanup() {
    rm -f ${tmpscript}
    exit
}

trap cleanup 0 1 2 5 15

# Does the script exist? If not, then exit immediatelly. In this case the
# Debian smartmontools package is not installed.
which update-smart-drivedb >/dev/null || exit 0

# Note, the script shipped with the Debian Wheezy smartmontools package
# contains a wrong URL. Because of that this must be modified to get the
# script working. To do not modify the original script, which is under
# APT control, we simply create a copy of the script, modify the URL and
# finally execute it.
cp -p $(which update-smart-drivedb) ${tmpscript}
sed -i "/^SRCEXPR/{s#=.*#='http://sourceforge.net/p/smartmontools/code/HEAD/tree/\$location/smartmontools/drivedb.h?format=raw'#}" ${tmpscript}
echo "Updating smartmontools drive database ..."
eval ${tmpscript}

exit 0
