# deb/openmediavault-clamav/etc/systemd/system/clamav-freshclam.service.d/openmediavault.conf
[Service]
ExecStartPre=/usr/bin/mkdir -p /var/lib/clamav/
ExecStartPre=/usr/bin/chmod 770 /var/lib/clamav/
ExecStartPre=/usr/bin/chown clamav:clamav /var/lib/clamav/
