# deb/openmediavault/etc/systemd/system/wsdd.service.d/openmediavault.conf
[Unit]
# If the unit bound to is stopped, this unit will be stopped too.
BindsTo=smbd.service
