# deb/openmediavault/etc/systemd/system/tmp.mount
# This file is part of openmediavault (https://www.openmediavault.org)
# WARNING: Do not edit this file.

[Unit]
Before=local-fs.target

[Mount]
What=tmpfs
Where=/tmp
Type=tmpfs
