# deb/openmediavault/debian/openmediavault.openmediavault-issue.service
[Unit]
Description=Generate the prelogin message
Wants=network-online.target
After=network-online.target
Before=getty.target

[Service]
Type=oneshot
ExecStart=/usr/bin/salt-call --local --retcode-passthrough --no-color state.orchestrate omv.deploy.issue
StandardOutput=null
StandardError=journal
RemainAfterExit=yes

[Install]
WantedBy=getty.target
